# src/apps/sysstat-6.0.1/sysstat.cron.hourly
#!/bin/sh
# Run system activity accounting tool every 10 minutes
/home/shearer/Projects/WAP4410N/build/../build/rootfs.build/lib/sa/sa1 -d 600 6 &
