# src/build/rootfs.build/etc/inittab
::sysinit:/etc/rc.d/rcS
::respawn:/sbin/getty ttyS0 115200
::shutdown:/etc/ath/killVAP all
::shutdown:/bin/umount -a
